# systemd/backup-nas.timer
[Unit]
Description=Run NAS backups daily.

[Timer]
OnBootSec=10s
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
